# rootfs/etc/services.d/cups-discovery/run
#!/bin/bash
set -euo pipefail
exec /opt/cups-autodiscover.sh
