# rootfs/etc/services.d/discovery-ui/run
#!/bin/bash
set -euo pipefail
# Discovery UI run: only run if enabled by a control file. Otherwise exit cleanly.
CONTROL_FILE="/run/cups/enable_discovery_ui"
if [ ! -f "$CONTROL_FILE" ]; then
	echo "[DISCOVERY-UI-SVC] Discovery UI disabled (control file $CONTROL_FILE not present). Exiting." >&2
	exit 0
fi
exec /opt/discovery-ui.sh

