# rootfs/etc/s6-overlay/s6-rc.d/init-cups/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

# Logging
log_info() { if command -v bashio >/dev/null 2>&1; then bashio::log.info "$*"; else echo "[init-cups] $*"; fi }
log_error() { if command -v bashio >/dev/null 2>&1; then bashio::log.error "$*"; else echo "[init-cups] ERROR: $*" >&2; fi }

log_info "Starting CUPS initialization..."

# Load configuration from Home Assistant
if command -v bashio >/dev/null 2>&1; then
    ADMIN_USER=$(bashio::config 'admin_username')
    ADMIN_PASS=$(bashio::config 'admin_password')
    SSL_ENABLED=$(bashio::config 'ssl_enabled')
    LOG_LEVEL=$(bashio::config 'log_level')
else
    info "Bashio not found, using defaults"
    ADMIN_USER="admin"
    ADMIN_PASS="admin"
    SSL_ENABLED="true"
    LOG_LEVEL="info"
fi

log_info "Configuration loaded: user=$ADMIN_USER, ssl=$SSL_ENABLED"

# Create directories
mkdir -p /var/log/cups /var/run/cups /var/cache/cups /var/spool/cups /etc/cups/ssl /run/cups /data/cups
chown -R lp:lp /var/log/cups /var/run/cups /var/cache/cups /var/spool/cups /etc/cups/ssl /run/cups

# Copy default config if /etc/cups is empty
if [[ ! -f /etc/cups/cupsd.conf ]]; then
    log_info "Copying default CUPS configuration..."
    cp -r /etc/cups.default/* /etc/cups/
fi

# Restore persistent config from /data (if it exists)
if [[ -d /data/cups/etc ]]; then
    log_info "Restoring persistent configuration from /data/cups..."
    cp -rf /data/cups/etc/* /etc/cups/ 2>/dev/null || true
fi

# Generate cupsd.conf via helper script
log_info "Generating cupsd.conf..."
/usr/local/bin/cups-config.sh \
    "$ADMIN_USER" \
    "$ADMIN_PASS" \
    "$SSL_ENABLED" \
    "$LOG_LEVEL"

# Generate SSL certificates
if [[ "$SSL_ENABLED" == "true" ]] && [[ ! -f /etc/cups/ssl/server.crt ]]; then
    log_info "Generating SSL certificates..."
    /usr/local/bin/generate-ssl.sh
fi

# Validate configuration
log_info "Validating CUPS configuration..."
if ! cupsd -t 2>/tmp/cupsd-test.log; then
    log_error "CUPS configuration is invalid:"
    cat /tmp/cupsd-test.log
    exit 1
fi

# Backup config to /data for persistence
mkdir -p /data/cups/etc
cp -rf /etc/cups/* /data/cups/etc/ 2>/dev/null || true

log_info "CUPS initialization completed successfully"
