# rootfs/etc/services.d/cups-ha-integration/run
#!/bin/bash
set -eu
# CUPS HA Integration / Simple Monitor

echo "[INFO] CUPS HA Integration service starting..."

# Minimal monitor to log CUPS status. If cupsd exits, log and exit with non-zero
# so HA can restart the add-on, rather than restarting cupsd inside the container.
while true; do
    if ! pgrep -x cupsd >/dev/null 2>&1; then
        echo "[ERROR] CUPS daemon not running; exiting to allow HA to restart the add-on"
        exit 1
    fi
    sleep 10
done
